# crontab
*/10 *  *   *   *     semaphore --max=1 --file=/tmp/update.semaphore -- cronlog  --logfile=/tmp/update.cronlog --timestamp perl /home/perldocjp/code/script/update.sh
